# server/systemd/scc-smart-poll.timer
# /etc/systemd/system/scc-smart-poll.timer — 每 30 分鐘收集一次硬碟 SMART（開機後 3 分鐘先跑一次）。
[Unit]
Description=SCC 硬碟 SMART 健康度收集排程

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
