# server/systemd/scc-attendance-pull@.timer
# scc-attendance-pull@.timer — 每 1 分鐘觸發一次 scc-attendance-pull@.service。
#
# v194：原每 10 分鐘改每 1 分鐘——實習生打卡後最壞要等 10 分鐘＋隨機延遲才進系統，使用者回報
# 體感太慢（打卡完重新整理一直看不到）。每輪成本僅一次 OAuth token refresh＋一次 Drive 下載，
# minutely（dev/prod 兩實例合計 ~2880 次/日）遠低於 Google API 配額，無量的顧慮。
#
# 安裝方式與參數說明見同目錄 README.md。
#
# ⚠ 修改本檔後需在 server 上重新安裝生效：
#   sudo cp server/systemd/scc-attendance-pull@.timer /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl restart scc-attendance-pull@dev.timer scc-attendance-pull@prod.timer

[Unit]
Description=每 1 分鐘拉取一次打卡紀錄（%i）

[Timer]
OnCalendar=*-*-* *:*:00
# 隨機延遲：避免 dev/prod 等多實例、或系統重開機後大量 timer 同時觸發時集中打 Google API。
RandomizedDelaySec=15
# 錯過的執行（例如機器重開機、timer 服務暫停期間）於下次啟動後補跑一次，
# 避免長時間斷線後累積過多筆待拉取紀錄。
Persistent=true

[Install]
WantedBy=timers.target
