# server/systemd/scc-gc-sync@.timer
# scc-gc-sync@.timer — 每 5 分鐘觸發一次 scc-gc-sync@.service。
#
# 對映 GAS setupGcSyncTrigger 的 every 5 min time trigger；非同步時段的 tick 由
# scripts/gc-sync-tick.js 內的時段閘自行跳過（exit 0），timer 不做時段判斷。
#
# 安裝方式與參數說明見同目錄 README.md。

[Unit]
Description=每 5 分鐘觸發一次 GC 日曆同步 tick（%i）

[Timer]
OnCalendar=*:0/5
# 隨機延遲：避免 dev/prod 等多實例、或系統重開機後大量 timer 同時觸發時集中打 Google API。
RandomizedDelaySec=30
# 同步是「比對現況」型工作，錯過的 tick 不需要補跑——下一個 5 分鐘的 tick 看到的就是最新現況；
# 與 pull 類 timer（Persistent=true 補拉累積資料）語意不同，故此處不設 Persistent。

[Install]
WantedBy=timers.target
