# server/systemd/scc-smart-poll.service
# /etc/systemd/system/scc-smart-poll.service
# 硬碟 SMART 健康度收集（系統全域單一實例——硬體為 dev/prod 共用，不用 @ 模板）。
# 以 root 執行（smartctl 讀 SMART 需要），輸出 /var/lib/scc-smart/smart.json（0644）供
# scc-dev/scc-prod 的 adminGetDiskHealth action 唯讀。app 程序不需任何提權。
#
# 安裝：
#   sudo cp scc-smart-poll.service scc-smart-poll.timer /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now scc-smart-poll.timer
#
# ⚠ 注意：ExecStart 指向 scc-prod checkout 內的腳本（隨 prod 部署自動更新）。此檔為
#   scc-s-admin 可寫、root 執行——在目前 scc-s-admin 本就有完整 sudo 的現實下無額外風險，
#   但日後若收斂 sudo 權限，應把腳本改為 root 擁有的固定路徑（如 /usr/local/sbin/）。
[Unit]
Description=SCC 硬碟 SMART 健康度收集（寫 /var/lib/scc-smart/smart.json）

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/node /home/scc-s-admin/scc-prod/server/scripts/smart-poll.js
